# com.openclaw.memory-api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.memory-api</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/jared/Projects/openclaw-memory/.venv/bin/python3</string>
        <string>/Users/jared/Projects/openclaw-memory/server.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/jared/Projects/openclaw-memory</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>MEMORY_DB</key>
        <string>/Users/jared/Projects/openclaw-memory/memory.db</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/memory-api.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/memory-api-err.log</string>
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
